# contrib/systemd/bluetooth-getty-session@.service
[Unit]
Description=Getty on RFCOMM tty %I
Documentation=man:agetty(8)
Documentation=https://github.com/sorah/bluetooth-getty
# Deliberately no BindsTo= / After=dev-%i.device. Lifecycle is owned
# by the bluetooth-getty daemon: it primes /dev/rfcommN via open()
# before calling StartUnit (so the device exists when agetty opens
# it), and it explicitly calls StopUnit in Profile1.RequestDisconnection
# when the peer disconnects. Relying on dev-%i.device was fragile
# because systemd's device unit is only created when udev tags the
# device with TAG+="systemd", which in practice was intermittent.
After=systemd-user-sessions.service
IgnoreOnIsolate=yes
Conflicts=rescue.service

[Service]
# session-proxy interposes a PTY between the rfcomm tty and
# agetty/login, so that util-linux login's unconditional vhangup()
# only affects the PTY slave — not the underlying RFCOMM DLC. Without
# this, vhangup() destroys the Bluetooth connection before login can
# reopen the tty. (shadow-utils login didn't call vhangup, but
# Ubuntu 26.04+ ships util-linux login.)
ExecStart=-/usr/sbin/bluetooth-getty session-proxy -- /sbin/agetty --keep-baud 115200,57600,38400,9600 - $TERM
Type=idle
# One login session per RFCOMM connection. After the user logs out
# the session-proxy exits and the daemon's SessionWatcher detects
# the unit going inactive. The peer must reconnect for a fresh
# session — reconnection yields a new NewConnection → fresh
# rfcomm_dev → fresh tty_struct.
Restart=no
UtmpIdentifier=%I
StandardInput=tty
StandardOutput=tty
TTYPath=/dev/%I
TTYReset=yes
# Deliberately NOT TTYVHangup=yes: for an rfcomm tty the vhangup()
# would trigger the kernel's hangup callback, destroying the RFCOMM
# DLC. Skipping vhangup is safe because the tty is freshly-created
# per connection with no stale state.
IgnoreSIGPIPE=no
SendSIGHUP=yes
ImportCredential=agetty.*
ImportCredential=login.*
